# stock-report-site/com.william.linecommands.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.william.linecommands</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/mavicair2tw/.openclaw/workspace/stock-report-site-repo/stock-report-site &amp;&amp; ./.venv/bin/python scripts/line_command_worker.py</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/mavicair2tw/.openclaw/workspace/stock-report-site/logs/line-commands.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mavicair2tw/.openclaw/workspace/stock-report-site/logs/line-commands.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
